# modules/wmem-outbox/install/wmem-outbox.service
[Unit]
Description=wmem-outbox local proxy daemon (wmem multi-instance PR-D)
Documentation=https://github.com/andrkelemen/wmem
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/node %h/wmem/modules/wmem-outbox/src/server.mjs
Restart=on-failure
RestartSec=5

# Defaults; override per-machine via Environment= or drop-in.
Environment=WMEM_UPSTREAM_HOST=127.0.0.1
Environment=WMEM_UPSTREAM_PORT=4200
Environment=WMEM_OUTBOX_PORT=4201
Environment=WMEM_OUTBOX_BIND=127.0.0.1
Environment=LOG_LEVEL=info

StandardOutput=journal
StandardError=journal
SyslogIdentifier=wmem-outbox

[Install]
WantedBy=default.target
